# home/kpanic/.config/autostart/ibus.desktop
[Desktop Entry]
Name=Ibus
Exec=ibus-daemon -rdx
Icon=redshift
Terminal=false
Type=Application
Categories=Utility;
StartupNotify=true
