# home/kpanic/.config/autostart/unclutter.desktop
[Desktop Entry]
Name=unclutter
GenericName=unclutte
Type=Application
Terminal=false
Exec="/usr/bin/unclutter"
